# OPOL/core/configs/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=uvicorn main:app --host 0.0.0.0 --port 8081
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:worker]
command=celery -A celery_worker.celery_app worker --loglevel=info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
